# debian/cron.d
# Launch the command to pull services from Nagios and create necessary tickets
# in Mantis
0 * * * * nagios2mantis_security.py
